# com.repcompanion.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.repcompanion.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/thomassoderberg/.gemini/antigravity/scratch/test/RepCompanion 2/start-server.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/thomassoderberg/.gemini/antigravity/scratch/test/RepCompanion 2</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/repcompanion-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/repcompanion-launchd-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_ENV</key>
        <string>development</string>
        <key>USE_AI_V3</key>
        <string>true</string>
        <key>USE_AI_V3_ONBOARDING</key>
        <string>true</string>
        <key>GEMINI_MODEL</key>
        <string>gemini-2.5-flash</string>
    </dict>
</dict>
</plist>

